# conf/upstart.conf
start on runlevel [2345]
stop on runlevel [016]

respawn

setgid wst
setuid wst

chdir /var/www/webstandardsdays.ru/twitter/www/src

exec node index.js
