# conf/upstart.conf
start on runlevel [2345]
stop on runlevel [016]

respawn

setuid www-data
setgid wst

chdir /var/www/webstandardsdays.ru/twitter/htdocs/src

exec node index.js
